# supervisor.conf
[program:board]
command = /opt/board/run_gunicorn.sh
user = board
stdout_logfile = /var/log/board/supervisor-gunicorn.log
redirect_stderr = true
